# cron/etc/cron.weekly/distfiles
#!/bin/sh
LOG=/var/log/cron.log
echo "*** $(date -R) ***" >> $LOG
eclean -Cd distfiles >> $LOG
