# cron/etc/cron.weekly/trim
#!/bin/sh
LOG=/var/log/cron.log
echo "*** $(date -R) ***" >> $LOG
fstrim -av >> $LOG
